# clem-schedule.service
[Unit]
Description=Clem's Walk Schedule Web Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Dev/clem_schedule
ExecStart=/usr/bin/python3 /home/pi/Dev/clem_schedule/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
